# image/ansible/coldstone/files/supervisor/dockerd.conf
[program:dockerd]
command=/usr/local/bin/dockerd
autostart=true
autorestart=false
stderr_logfile=/var/log/dockerd.err.log
stdout_logfile=/var/log/dockerd.out.log